# deploy/systemd/toxiflow-backup.timer
[Unit]
Description=Daily backup for ToxiFlow

[Timer]
OnCalendar=*-*-* 03:15:00
Persistent=true
Unit=toxiflow-backup.service

[Install]
WantedBy=timers.target
